# external_deploy/hf_space/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:smartbot]
command=/app/start.sh
autostart=true
autorestart=true
startretries=999
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"

[program:healthcheck]
command=bash -c "while true; do curl -sf http://localhost:7860/ || (echo 'Bot dead, restarting...'; supervisorctl restart smartbot); sleep 60; done"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
